# supervisor.ini
[program:server]
command = python manage.py server
directory = /mypath
usr= root
autostart = true
autorestart = true
startsecs = 5
stdout_logfile=/var/log/eclogue-out.log
stderr_logfile=/var/log/eclogue-error.log

[program:worker]
command =python manage.py worker
directory = /mypath
usr= root
autostart = true
autorestart = true
startsecs = 5
stdout_logfile=/var/log/worker-out.log
stderr_logfile=/var/log/worker-error.log
[group:eclogue]
programs=server,worker
