# supervisor.ini
[program:eclogue]
command =/usr/local/bin/pipenv run python manage.py start
direcotry = %(here)s
usr= root
autostart = true
autorestart = true
startsecs = 5
stdout_logfile=/var/log/eclogue-out.log
stderr_logfile=/var/log/eclogue-error.log

[program:eclogue]
command =/usr/local/bin/pipenv run python manage.py worker
usr= root
autostart = true
autorestart = true
startsecs = 5
stdout_logfile=/var/log/worker-out.log
stderr_logfile=/var/log/worker-error.log
